# ops/sinitaivas-live.service
[Unit]
Description=Sinitaivas Live Streamer
After=network-online.target
Wants=network-online.target

[Service]
User=mehrdad
Group=mehrdad
WorkingDirectory=/home/mehrdad/projects/bluesky_scrapper/sinitaivas-live
Environment=PYTHONUNBUFFERED=1
# call the venv’s interpreter directly — no shell, no pipenv
ExecStart=/home/mehrdad/projects/bluesky_scrapper/sinitaivas-live/.venv/bin/python -m sinitaivas_live.main --mode resume
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sinitaivas-live

[Install]
WantedBy=multi-user.target